# src/networkd_connectivity/_resources/systemd/agentx-connectivity.service
[Unit]
Description=SNMP AgentX sub-agent for networkd-connectivity
After=snmpd.service
Wants=snmpd.service

[Service]
Type=notify
ExecStart=/usr/bin/agentx-connectivity
Environment=AGENTX_SOCKET=/var/agentx/master

[Install]
WantedBy=multi-user.target
